# deploy/systemd/phenopi-scheduler.service.in
[Unit]
Description=Phenopi experiment scheduler
After=local-fs.target time-sync.target
Wants=time-sync.target

[Service]
Type=simple
User=@PHENOPI_USER@
Group=@PHENOPI_GROUP@
EnvironmentFile=/etc/phenopi/phenopi.env
WorkingDirectory=@PHENOPI_ROOT@
ExecStart="@PHENOPI_PYTHON@" -m scripts.scheduling.scheduler
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths="@PHENOPI_RUNTIME_DIR@" "@PHENOPI_CAPTURE_DIR@"

[Install]
WantedBy=multi-user.target
