# deploy/systemd/phenopi-gui.service.in
[Unit]
Description=Phenopi web interface
After=network-online.target phenopi-scheduler.service
Wants=network-online.target

[Service]
Type=simple
User=@PHENOPI_USER@
Group=@PHENOPI_GROUP@
EnvironmentFile=/etc/phenopi/phenopi.env
WorkingDirectory=@PHENOPI_ROOT@
ExecStart="@PHENOPI_PYTHON@" -m gui.app
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths="@PHENOPI_RUNTIME_DIR@" "@PHENOPI_CAPTURE_DIR@"

[Install]
WantedBy=multi-user.target
